# {{cookiecutter.project_slug}}/etc/crontab
* * * * * cd /home/app/ && /usr/local/bin/pipenv run python manage.py build_index >/home/app/cron.log 2>&1
